# roles/restic/files/restic.timer
[Unit]
Description=Run restic backup

[Timer]
OnCalendar=*-*-* 04:00:00
Persistent=true
RandomizedDelaySec=1h

[Install]
WantedBy=timers.target
